# systemd/user/shuffle_wallpaper.timer
[Unit]
Description="Shuffle Wallpapers"
[Timer]
OnCalendar=*:0/15
Unit=shuffle_wallpaper.service
[Install]
WantedBy=multi-user.target
