# systemd/user/shuffle_wallpaper.service
[Unit]
Description="Shuffle Wallpaper Script"
[Service]
ExecStart=/home/xempt/.scripts/shuffle_wallpaper.sh
